# supervisor/cbms.conf
[program:cbms]
user = root 
command = /usr/cbms/bin/svm --daemon
autostart = true
autorestart = true
startsecs = 5

