# supervisor/nginx.conf
[program:nginx]
user = nginx
command = /usr/sbin/nginx 
autostart = true
autorestart = true
startsecs = 5

